# .config/systemd/user/blueman-applet.service
[Service]
Type=simple
ExecStart=blueman-applet

[Install]
WantedBy=default.target

