# .config/systemd/user/hyprpaper.service
[Service]
Type=simple
ExecStart=hyprpaper

[Install]
WantedBy=default.target

