# .config/systemd/user/kwallet.service
[Service]
Type=simple
ExecStart=kwalletd6

[Install]
WantedBy=default.target

