# .config/systemd/user/kwallet.service
[Unit]
After=wayland-session@hyprland.desktop.target
Requires=wayland-session@hyprland.desktop.target

[Service]
Type=simple
ExecStart=kwalletd6

[Install]
WantedBy=default.target

