# .config/systemd/user/wl-clip.service
[Service]
Type=simple
ExecStart=wl-clip-persist --clipboard both

[Install]
WantedBy=default.target

