# .config/systemd/user/wl-paste@.service
[Service]
Type=simple
ExecStart=wl-paste --type %I --watch cliphist store

[Install]
WantedBy=default.target

